# systemd任务/etcd定时备份/etcd3_staas.timer
[Unit]
Description=Run etcd3 staas
Documentation=

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=multi-user.target